# misc/xdg-desktop-portal-luminous.service.in
[Unit]
Description=Xdg Desktop Portal backend for wlroots-based compositors written with rust
PartOf=graphical-session.target
After=graphical-session.target
ConditionEnvironment=WAYLAND_DISPLAY

[Service]
Type=dbus
BusName=org.freedesktop.impl.portal.desktop.luminous
ExecStart=@xdg_install_dir@/xdg-desktop-portal-luminous
Restart=on-failure
